# supervisor-tempmonitor.conf
[program:app-tempmonitor-web]
command = tempmonitor web

[program:app-tempmonitord]
command = tempmonitor daemon